# myscript.service
[Unit]
Desription=My Script
After=network.target

[Service]
ExecStart=/etc/ntp.sh

[Install]
WantedBy=default.target
